# scripts/AGL/layers/meta-mylayer/recipes-raspberry/utils-scripts/files/CarControl.service
[Unit]
Description=CarControl Application

[Service]
Type=simple
ExecStart=/usr/bin/CarControl

[Install]
WantedBy=multi-user.target
